# build/service/coldfusion2016/run
#!/bin/bash

set -e

exec /opt/cf2016/cfusion/bin/coldfusion start
